# files/etc/systemd/system/restic-server.timer
[Unit]
Description=Manage restic rest-server

[Timer]
OnCalendar=*-*-* 21:05:00

[Install]
WantedBy=timers.target

